# deploy/screen-export.service
[Unit]
Description=Screen-Export FastAPI Backend
Documentation=https://github.com/Achinoise1/Screen-Export
After=network.target

[Service]
Type=simple
# ── Adjust these two lines to match your deploy server ──────────────────────
User=www-data
WorkingDirectory=/data/projs/Screen-Export
# ────────────────────────────────────────────────────────────────────────────

ExecStart=/data/projs/Screen-Export/.venv/bin/python backend/main.py

# Sub-path configuration — must match the nginx location block
Environment="BACKEND_ROOT_PATH=/tools/screenshot-export"
Environment="BACKEND_PORT=8001"

# Optional: relocate data outside the project directory
Environment="SCREEN_EXPORT_DATA_DIR=/var/screen-export-data"

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=screen-export

[Install]
WantedBy=multi-user.target
